# systemd/claw-face-evening.timer
[Unit]
Description=Claw Face evening routine at 22:00

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
